# salon-core/scripts/systemd/salon-core-partitions.timer
[Unit]
Description=salon-core monthly partitions timer

[Timer]
OnCalendar=*-*-01 04:10:00
Persistent=true

[Install]
WantedBy=timers.target
