# salon-core/scripts/systemd/salon-core-retention.service
[Unit]
Description=salon-core retention cleanup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/salon-core/env
Environment=RUN_RETENTION=1
Environment=RUN_PARTITIONS=0
ExecStart=/opt/salon-core/scripts/cron/maintenance.sh
